# build/config/includes.chroot/etc/systemd/system/wrathos-firstboot.service
[Unit]
Description=WrathOS First Boot Setup
After=graphical-session.target
ConditionPathExists=!/var/lib/wrathos-firstboot-done

[Service]
Type=oneshot
ExecStart=/usr/bin/wrathos-firstboot.sh
RemainAfterExit=yes

[Install]
WantedBy=graphical-session.target
