# build/config/includes.chroot/etc/systemd/system/wrathos-firstboot.service
[Unit]
Description=WrathOS First Boot Setup
After=multi-user.target
ConditionPathExists=!/var/lib/wrathos-firstboot-done

[Service]
Type=oneshot
ExecStart=/usr/bin/wrathos-firstboot-system.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
